# supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=/usr/lib/postgresql/9.1/bin/postgres -D /var/lib/postgresql/9.1/flexibee -c config_file=/etc/postgresql/9.1/flexibee/postgresql.conf
user=postgres
autorestart=true

[program:postgresqlMain]
command=/usr/lib/postgresql/9.1/bin/postgres -D /var/lib/postgresql/9.1/main -c config_file=/etc/postgresql/9.1/main/postgresql.conf
user=postgres
autorestart=true

[program:flexibee]
command=/supervisord-flexibee.sh
#user=winstrom
#numprocs=1
#autostart=true
#autorestart=true